# install/supervisord/conf.d/kartoza-geoserver.conf
# Use supervisor to restart docker containers, specifically for when services it depends on is not ready
# https://coderwall.com/p/aqcbzw/server-deployments-on-ubuntu-with-docker-and-supervisor

[program:kartoza-geoserver]
priority=2
autostart=true
startsecs=11
startretries=10
autorestart=true
stopsignal=TERM
stopwaitsecs=60
command=bash %(here)s/conf.d/kartoza-geoserver.sh


